# hosts/cadmium/crontab
# m h dom mon dow       command
4-59/5 * * * *		adt-play/cronjob
22 0 * * Mon		savelog -n7 adt-play/overall.log
58 7 * * *		adt-play/cleanup
